# bandles/p-board/raspbian/home/pi/crontab
MAILTO=""
@reboot sudo ~/u.board/cront/gpio.sh
@reboot /usr/bin/chromium-browser --kiosk http://localhost/u.board


