# cron/crontab
MAILTO=""

* * * * * /home/pi/u.board/cron/flag.sh
