# raspbian/bandles/a-board/crontab
MAILTO=""

@reboot /usr/bin/python ~/u.board/raspbian/bandles/a-board/bin/start.py
